# sandbox/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:xvfb]
command=Xvfb :0 -screen 0 1280x720x24 -nolisten tcp -ac
priority=10
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:x11vnc]
command=x11vnc -display :0 -forever -shared -localhost -rfbport 5900 -nopw
priority=20
startsecs=2
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:websockify]
command=websockify --web /usr/share/novnc 127.0.0.1:8080 127.0.0.1:5900
priority=30
startsecs=2
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:chromium]
command=/usr/bin/chromium --user-data-dir=/data --no-first-run --disable-session-crashed-bubble --window-position=0,0 --window-size=1280,720 --remote-debugging-address=127.0.0.1 --remote-debugging-port=9222 --proxy-server=http://127.0.0.1:8118 --no-sandbox --disable-dev-shm-usage --disable-extensions-except=/opt/agent-browser/chrome-extension --load-extension=/opt/agent-browser/chrome-extension https://example.com/
priority=40
startsecs=3
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:agent-browser-init]
command=/bin/sh -c 'sleep 5; agent-browser-cli restart'
priority=50
startsecs=0
autorestart=false
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:tor]
command=/usr/bin/tor --RunAsDaemon 0 --SocksPort 127.0.0.1:9050 --DataDirectory /tmp/tor
priority=60
startsecs=3
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:sandbox-proxy]
command=/usr/local/bin/sandbox-proxy
priority=70
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
